# systemd/user/pomodoro.service
[Unit]
Description=Manually trigger session idle state using swayidle

[Service]
ExecStart=/usr/bin/pkill -USR1 swayidle
