# systemd/user/pomodoro.timer
[Unit]
Description=Trigger session idle state after 25 minutes using swayidle

[Timer]
OnActiveSec=25min
RemainAfterElapse=no
